# dotfiles/config/zsh/profile.d/garuda-seamus.zsh
# Connection helpers
alias lxcL='ssh -t seamus-lxc tmux -u new-session -A -s lxc-L'
alias lxcR='ssh -t seamus-lxc tmux -u new-session -A -s lxc-R'

# Garuda-home vars
export MPD_HOST="/home/seamus/.mpd/socket"
export _JAVA_OPTIONS="-Dawt.useSystemAAFontSettings=on -Dswing.aatext=true"
