# dotfiles/config/zsh/profile.d/seamus-lxc.zsh
# Work helpers
alias ssh="signmykey -e; ssh"
alias scp="signmykey -e; scp"
alias tmuxp="signmykey -e; tmuxp"
alias ansible="signmykey -e; ansible"
alias ansible-playbook="signmykey -e; ansible-playbook"
alias htop='gotop'

# Work vars
export VAULT_ADDR=https://vault.maxim.services:8200
