# dotfiles/config/zsh/profile.d/seamus-pad.zsh
# source "main" profile
source $ZCONFDIR/profile.d/garuda-seamus.zsh
